# debian/oyo-portable-system-creator/usr/share/oyo-portable-system-creator/src/templates/firstboot.service
[Unit]
Description=OYO Portable Firstboot
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/oyo-firstboot

[Install]
WantedBy=multi-user.target
